# supervisord.conf
[supervisord]
nodaemon=true

[program:canal_flask]
command=python /bin/canal-flask.py
